# com.tahoesno.tmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Installs to ~/Library/LaunchAgents/com.tahoesno.tmux.plist
     Before loading: unload the old server plist if still active:
       launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.tahoesno.server.plist
     Install + start:
       cp com.tahoesno.tmux.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.tahoesno.tmux.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tahoesno.tmux</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/cortex/projects/tahoe-sno/tmux-serve.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>8788</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/tahoesno-tmux.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tahoesno-tmux.log</string>
</dict>
</plist>
